# com.caps.dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Job label - must be unique on your Mac -->
    <key>Label</key>
    <string>com.caps.dashboard</string>

    <!-- The script to run -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>REPLACE_WITH_FULL_PATH/update_dashboard.py</string>
    </array>

    <!-- Run at 9:00 AM every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Working directory (same folder as script) -->
    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_FULL_PATH</string>

    <!-- Log output for debugging -->
    <key>StandardOutPath</key>
    <string>REPLACE_WITH_FULL_PATH/logs/dashboard_run.log</string>

    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_FULL_PATH/logs/dashboard_error.log</string>

    <!-- Run even if Mac was asleep at 9 AM -->
    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
